# supervisord-syslog.conf

[supervisord]
nodaemon=true

[program:fluentd]
command=bash /startfluentsyslog.sh
